# iso-overlay/etc/rc.local
#

#Set xorg.conf
sh /root/bin/xconf.sh
sleep 4;
/usr/local/etc/init.d/slim start
